# rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon = true

[program:php-fpm]
command = php-fpm7.4 --allow-to-run-as-root --nodaemonize --fpm-config /etc/php/7.4/fpm/php-fpm.conf

[program:nginx]
command = nginx
